# openwrt/etc/hotplug.d/iface/12-cake-autorate
#!/bin/sh

[ "$ACTION" = ifup ] && /etc/init.d/sqm enabled && /etc/init.d/cake-autorate restart

